# scripts/blog-local/com.chifbay.blog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chifbay.blog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/futurx/Downloads/Chifbay/dist/scripts/blog-local/run-blog.sh</string>
  </array>
  <!-- Runs every day at 09:00 local time. To change the days, replace this
       <dict> with an <array> of <dict>s (each with Weekday 0=Sun..6=Sat + Hour). -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/futurx/Downloads/Chifbay/dist/scripts/blog-local/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/futurx/Downloads/Chifbay/dist/scripts/blog-local/launchd.err.log</string>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
